# contrib/asnumber.socket
[Unit]
Description=asnumber socket

[Socket]
ListenStream=/run/gunicorn/asnumer.socket

[Install]
WantedBy=sockets.target
